# scripts/treason.conf
description "Treason node.js server"
author      "chris"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/treason/deploy"

    exec sudo -u treason /usr/bin/nodejs /home/treason/deploy/server.js --port 80 >> /var/log/treason 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
